# addons/init/runit/kula/run
#!/bin/sh
# Kula - runit service

# Runtime dir for an optional Unix socket under /run/kula (web.unix_socket).
# /run is a tmpfs wiped each boot and kula does not create the socket's parent
# dir itself. Unused in the default TCP setup; harmless if so.
mkdir -p /run/kula
chown kula:kula /run/kula
chmod 0750 /run/kula

exec chpst -u kula:kula /usr/bin/kula --config /etc/kula/config.yaml serve 2>&1
